# site/fastd/templates/fastd.upstart.erb
# this file is managed by puppet

description "fastd"
author "Hubert Mühlhans <muehlhans.hubert@ekodia.de>"

start on started networking
stop on stopped networking

env CONFIG_FILE=<%= scope['::fastd::params::community_folder'] -%>/fastd.conf

pre-start script
  if [ ! -f /etc/default/fastd ]; then
    echo "/etc/default/fastd does not exist! - Aborting..."
    stop
    exit 0
  fi
  . /etc/default/fastd

  if [ ! -f $CONFIG_FILE ]; then
    echo "$CONFIG_FILE does not exist! - Aborting..."
    echo "Please create and configure $CONFIG_FILE to fix the problem."
    stop
    exit 0
  fi

end script

respawn
expect daemon

script
  . /etc/default/fastd

  /usr/bin/fastd --syslog-level info --config $CONFIG_FILE
end script

